# reallysimplebackup.cron
#
# cron-jobs for backup-rsync and backup-rotate
#

MAILTO=root

# every 60 minutes, sync backup (will automatically rotate backups as needed)
0 * * * *     root if [ -x /usr/bin/reallysimplebackup-rsync ]; then /usr/bin/reallysimplebackup-rsync >> /var/log/reallysimplebackup.log 2>&1; fi

# at every reboot, force backup-rsync
@reboot root if [ -x /usr/bin/reallysimplebackup-rsync ]; then /usr/bin/reallysimplebackup-rsync after-reboot >> /var/log/reallysimplebackup.log 2>&1; fi

